# crontab
# Run scraper every day at 2:00
0 2 * * * /usr/local/bin/python /scraper/main.py >> /var/log/cron.log 2>&1
